# cron.d/localmon
PATH=/bin:/usr/bin:/sbin:/usr/sbin
SHELL=/bin/bash

* * * * * root /opt/localmon/usr/local/bin/compareCommandOutput.pl 'zpool status' compare normal 
* * * * * root for i in {1..30}; do /opt/localmon/usr/local/bin/compareCommandOutput.pl lsblk save normal 30000; sleep 2; done
* * * * * root for i in {1..30}; do /opt/localmon/usr/local/bin/compareCommandOutput.pl lsusb save normal 30000; sleep 2; done
