# cron.d/localmon
PATH=/bin:/usr/bin:/sbin:/usr/sbin
SHELL=/bin/bash

0 * * * * root /opt/localmon/bin/compareCommandOutput.pl '/opt/localmon/bin/zpool_status_localmon' compare normal 3600000
* * * * * root for i in {1..30}; do /usr/local/bin/compareCommandOutput.pl lsblk save normal 30000; sleep 2; done
* * * * * root for i in {1..30}; do /usr/local/bin/compareCommandOutput.pl lsusb save normal 30000; sleep 2; done
